# etc/supervisor.d/confd.ini
[program:confd]
command=/usr/bin/confd -node %(ETCDHOST)s:%(ETCDPORT)s
