# etc/supervisor.d/confd.ini
[program:confd]
command=/usr/bin/confd -node %(ENV_ETCDHOST)s:%(ENV_ETCDPORT)s
